# roles/redis/templates/redis@.service.j2
# This file is automatically generated - any manual changes will be overwritten.
# {{ ansible_managed | replace(";", "\n#") }}

[Unit]
Description=Advanced key-value store at %i
After=network.target

[Service]
Type=notify
User=redis
Group=redis
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadOnlyPaths=/
LimitNOFILE=65535
RuntimeDirectory=redis
RuntimeDirectoryMode=2755
UMask=007
ReadWritePaths=-/var/lib/redis
ReadWritePaths=-/var/log/redis
ReadWritePaths=-/run/redis
PIDFile=/run/redis/%i.pid
ExecStartPre=/usr/bin/test -f /etc/redis/%i.conf
ExecStart=/usr/bin/redis-server /etc/redis/%i.conf --daemonize yes --supervised systemd
Restart=on-failure
TimeoutStartSec=5s
TimeoutStopSec=5s

[Install]
WantedBy=multi-user.target
