# roles/rabbitmq/templates/epmd/system/epmd.socket.j2
# This file is automatically generated - any manual changes will be overwritten.
# {{ ansible_managed | replace(";", "\n#") }}

[Unit]
Description=Erlang Port Mapper Daemon Activation Socket

[Socket]
ListenStream=4369
BindIPv6Only=both
Accept=no

[Install]
WantedBy=sockets.target
